# ops/systemd/slop-publish.timer
[Unit]
Description=Trigger a Slop University publish run every 6 hours

[Timer]
OnCalendar=*-*-* 00/6:00
RandomizedDelaySec=15min
Persistent=true
Unit=slop-publish.service

[Install]
WantedBy=timers.target
