# ops/systemd/slop-publish.service
[Unit]
Description=Slop University: autonomous publish run (scan, generate, commit; wrapper validates and pushes)
Documentation=file:///home/ben/projects/slop-university/ops/cron-publish.sh
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/home/ben/projects/slop-university
ExecStart=/home/ben/projects/slop-university/ops/cron-publish.sh
TimeoutStartSec=2h
