# ops/systemd/slop-publish.timer
[Unit]
Description=Trigger a Slop University publish run every hour

[Timer]
OnCalendar=*-*-* *:00
RandomizedDelaySec=5min
Persistent=true
Unit=slop-publish.service

[Install]
WantedBy=timers.target
